# .config/systemd/user/xcompmgr.service
[Unit]
Description=X Compositing

[Service]
ExecStart=/usr/bin/xcompmgr

[Install]
WantedBy=graphical.target
